# services/sync/deploy/agent-sync.service
[Unit]
Description=Whispers Agent Sync Service (FastAPI Dashboard)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/whispers/WhispersoftheSevenKingdoms/services/sync
Environment=GITHUB_WEBHOOK_SECRET=CHANGE_ME
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/whispers/WhispersoftheSevenKingdoms/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
